# scripts/symharix.service
[Unit]
Description=SymHarix Telegram-first AI Supervisor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SYMHARIX_USER__
WorkingDirectory=__SYMHARIX_WORKDIR__
Environment="PATH=__SYMHARIX_PATH__"
Environment="PORT=3000"
ExecStart=__SYMHARIX_BUN__ --env-file=__SYMHARIX_WORKDIR__/.env run __SYMHARIX_WORKDIR__/src/cli/index.ts --port 3000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
