# scripts/symphonyness.service
[Unit]
Description=symphonyness Harnessed Code Development Platform
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=/home/%USER%/symphonyness
Environment=PATH=/usr/bin:/usr/local/bin:/home/%USER%/.bun/bin
ExecStart=/usr/bin/env bun run start
Restart=always
RestartSec=10
StandardOutput=append:/home/%USER%/symphonyness/symphonyness.log
StandardError=append:/home/%USER%/symphonyness/symphonyness.log

[Install]
WantedBy=multi-user.target
